# packaging/DEBIAN/postinst
#!/bin/sh
set -e

# Create user if it doesn't exist
if ! id "sysupdate-api" >/dev/null 2>&1; then
    useradd -r -s /bin/false sysupdate-api
fi

# Set permissions
chown -R sysupdate-api:sysupdate-api /opt/pinsdaemon

# Ensure sudoers file is owned by root (Critical for sudo to work)
chown root:root /etc/sudoers.d/sysupdate-api
chmod 440 /etc/sudoers.d/sysupdate-api

# Ensure upgrade script is owned by root (Security)
chown root:root /usr/local/bin/system-upgrade.sh
chmod 755 /usr/local/bin/system-upgrade.sh
chown root:root /usr/local/bin/manage-samba.sh
chmod 755 /usr/local/bin/manage-samba.sh

# Ensure service file is owned by root
chown root:root /etc/systemd/system/sysupdate-api.service
chmod 644 /etc/systemd/system/sysupdate-api.service

# Setup Python virtual environment
cd /opt/pinsdaemon
if [ ! -d "venv" ]; then
    echo "Creating virtual environment..."
    python3 -m venv venv
fi

# Install dependencies
if [ -f "requirements.txt" ]; then
    echo "Installing Python dependencies..."
    ./venv/bin/pip install -r requirements.txt
fi

# Fix ownership of venv created by root
chown -R sysupdate-api:sysupdate-api /opt/pinsdaemon/venv

# Systemd setup
systemctl daemon-reload
systemctl enable sysupdate-api
systemctl restart sysupdate-api

exit 0
