# packaging/DEBIAN/postinst
#!/bin/sh
set -e

# Create user if it doesn't exist
if ! id "sysupdate-api" >/dev/null 2>&1; then
    useradd -r -s /bin/false sysupdate-api
fi

# Set permissions
chown -R sysupdate-api:sysupdate-api /opt/pinsdaemon

# Ensure sudoers file is owned by root (Critical for sudo to work)
chown root:root /etc/sudoers.d/sysupdate-api
chmod 440 /etc/sudoers.d/sysupdate-api

# Ensure upgrade script is owned by root (Security)
if [ -f "/usr/local/bin/system-upgrade.sh" ]; then
    chown root:root /usr/local/bin/system-upgrade.sh
    chmod 755 /usr/local/bin/system-upgrade.sh
fi

if [ -f "/usr/local/bin/wifi-connect.sh" ]; then
    chmod 755 /usr/local/bin/wifi-connect.sh
fi

if [ -f "/usr/local/bin/wifi-automanage.py" ]; then
    chmod 755 /usr/local/bin/wifi-automanage.py
fi

if [ -f "/usr/local/bin/manage-samba.sh" ]; then
    chown root:root /usr/local/bin/manage-samba.sh
    chmod 755 /usr/local/bin/manage-samba.sh
fi

if [ -f "/usr/local/bin/wifi-connect.sh" ]; then
    chown root:root /usr/local/bin/wifi-connect.sh
    chmod 755 /usr/local/bin/wifi-connect.sh
fi

if [ -f "/usr/local/bin/wifi-automanage.py" ]; then
    chown root:root /usr/local/bin/wifi-automanage.py
    chmod 755 /usr/local/bin/wifi-automanage.py
fi

if [ -f "/home/pi/pins-scripts/hotspot.sh" ]; then
    chown pi:pi /home/pi/pins-scripts/hotspot.sh
    chmod 755 /home/pi/pins-scripts/hotspot.sh
fi

if [ -f "/usr/local/bin/wifi-scan.py" ]; then
    chown root:root /usr/local/bin/wifi-scan.py
    chmod 755 /usr/local/bin/wifi-scan.py
fi

# Ensure service file is owned by root
chown root:root /etc/systemd/system/sysupdate-api.service
chmod 644 /etc/systemd/system/sysupdate-api.service

# Setup Python virtual environment
cd /opt/pinsdaemon
if [ ! -d "venv" ]; then
    echo "Creating virtual environment..."
    python3 -m venv venv
fi

# Install dependencies
if [ -f "requirements.txt" ]; then
    echo "Installing Python dependencies..."
    ./venv/bin/pip install -r requirements.txt
fi

# Fix ownership of venv created by root
chown -R sysupdate-api:sysupdate-api /opt/pinsdaemon/venv

# Systemd setup
systemctl daemon-reload
systemctl enable sysupdate-api
systemctl restart sysupdate-api

exit 0
